# build_files/astroimmutable-firstboot.service
[Unit]
Description=AstroImmutable first login setup
ConditionPathExists=!%h/.local/state/astroimmutable/firstlogin.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/astroimmutable/firstlogin-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
